# src/subscript/rc.local
#!/bin/bash

./root/scripts/src/subscript/port8080_to_80_redirect.sh
